# deploy/mcp/odoo-mcp.service
[Unit]
Description=Odoo MCP Server (HTTP)
After=network.target odoo.service
Wants=odoo.service

[Service]
Type=simple
User=odoo
Group=odoo
WorkingDirectory=/opt/odoo/mcp-server
# Credentials in a restricted file (chmod 600, owned by odoo)
EnvironmentFile=/etc/odoo-mcp/credentials
ExecStart=/root/.local/bin/uv run --with mcp python3 /opt/odoo/mcp-server/server.py --http --host 127.0.0.1 --port 8200
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
