# deploy/mcp/odoo-mcp.service
[Unit]
Description=Odoo MCP Server (HTTP)
After=network.target odoo.service
Wants=odoo.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/odoo/mcp-server
# Credentials in a restricted file (chmod 600, owned by root)
EnvironmentFile=/etc/odoo-mcp/credentials
ExecStart=/root/.local/bin/uv run --with mcp python3 /opt/odoo/mcp-server/server.py --http --host 0.0.0.0 --port 8200
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
